# user/env.target
[Unit]
Description=Environment Target
Wants=dbus.socket
Requires=dbus.socket

[Install]
WantedBy=wm.target
