# user/env.target
[Unit]
Description=Environment Target
Requires=dbus.socket
Before=mystuff.target

[Install]
WantedBy=mystuff.target
