# user/mystuff.target
[Unit]
Description=Xinitrc Stuff
After=env.target

[Install]
Alias=default.target
WantedBy=wm.target
