# user/urxvtd-user.service
[Unit]
Description=Urxvt Daemon
After=env.target

[Service]
ExecStart=/usr/bin/urxvtd

[Install]
WantedBy=mystuff.target
